# rootfs/etc/services.d/nginx-entrypoint/run
#!/command/with-contenv bash

exec /docker-entrypoint.sh nginx -g "daemon off;"